# files/conf/logrotate/platform.cron
*/1 * * * * /usr/sbin/logrotate /etc/logrotate.d/platform"
