# files/conf/logrotate/platform.cron
*/1 * * * * /usr/sbin/logrotate /etc/logrotate.d/platform && logger -i -t "crontab" -p local0.notice "logrotate /etc/logrotate.d/platform"
